# systemd/system/pleiades-nexus.service
[Unit]
Description=Pleiades Nexus Omniversal
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/containment_controller
Restart=always
RestartSec=10
LimitNOFILE=4096
MemoryMax=1G
CPUQuota=100%

[Install]
WantedBy=multi-user.target
